# sparse/etc/systemd/system/sensorfwd.service.d/90-boot-order.conf
[Unit]
Requires=systemd-modules-load.service
ConditionPathExists=/dev/iio:device1
