# sparse/etc/systemd/system/sensorfwd.service.d/90-boot-order.conf
[Unit]
Requires=sensorfwd-iio-start.path
